# private_dot_config/private_systemd/user/foot-server@.socket.d/override.conf
[Unit]
After=graphical-session.target
BindsTo=graphical-session.target

[Install]
WantedBy=
WantedBy=graphical-session.target