# private_dot_config/private_systemd/user/rc.service
[Unit]
Description=Run commands

[Service]
ExecStart=/bin/sh -c 'bash $HOME/.config/rc.d/*'

[Install]
WantedBy=default.target
